# linux/services/example.timer
[Unit]
Description=Run my task every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=multi-user.target
